# supervisord.conf
[supervisord]
nodaemon=true
username=www-data

[program:php-fpm]
command=/usr/sbin/php-fpm8.2 -F
username=www-data
autostart=true
autorestart=true