# supervisord.conf
[supervisord]
nodaemon=true
username=www-data

[program:fpm-version]
command=/usr/sbin/fpm-version -F
username=www-data
autostart=true
autorestart=true