# supervisord.conf
[supervisord]
nodaemon=true
username=www-data

[program:php5.6-fpm]
command=/usr/sbin/php-fpm5.6 -F
username=www-data
autostart=true
autorestart=true