# supervisord.conf
[supervisord]
nodaemon=true
username=www-data

[program:php7.3-fpm]
command=/usr/sbin/php-fpm7.3 -F
username=www-data
autostart=true
autorestart=true