# supervisord.conf
[supervisord]
nodaemon=true
username=www-data

[program:php7.0-fpm]
command=/usr/sbin/php-fpm7.0 -F
username=www-data
autostart=true
autorestart=true