# supervisord.conf
[supervisord]
nodaemon=true
username=www-data

[program:php7.4-fpm]
command=/usr/sbin/php-fpm7.4 -F
username=www-data
autostart=true
autorestart=true