# supervisord.conf
[supervisord]
nodaemon=true

[program:php5.6-fpm]
command=/usr/sbin/php-fpm5.6 -F
autostart=true
autorestart=true