# supervisord.conf
[supervisord]
nodaemon=true

[program:php7.0-fpm]
command=/usr/sbin/php-fpm7.0 -F
autostart=true
autorestart=true