# supervisord.conf
[supervisord]
nodaemon=true
username=www-data

[program:php7.2-fpm]
command=/usr/sbin/php-fpm7.2 -F
username=www-data
autostart=true
autorestart=true